# owl-net-timer.timer
[Unit]
Description=Run Owl Network Timer

[Timer]
OnBootSec=15sec
OnUnitActiveSec=20sec

[Install]
WantedBy=timers.target