# owl-net-timer.service
[Unit]
Description=Owl Network Timer

[Service]
ExecStart=/usr/bin/python3 /etc/owl/network_manager.py