# owl-net-timer.timer
[Unit]
Description=Run Owl Network Timer

[Timer]
OnBootSec=10sec
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target